# main/centos7/docker/collect_server/log_collect/shell/crontab
*/1 * * * * /home/log_collect/log_collect/log_collect.sh
00 10 * * * /home/log_collect/countup.sh
40 10 * * * /home/log_collect/reset.sh