# main/centos7/docker/collect_server/log_collect/crontab
*/1 * * * * /root/src2/log_collect/log_collect.sh
00 10 * * * /root/src2/log_collect/countup.sh
40 10 * * * /root/src2/log_collect/reset.sh